# config/launchd/com.hexa-lang.sleep_mode.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.hexa-lang.sleep_mode — user 취침 중 주기 meta pass.

  Spec: $HEXA_LANG/tool/sleep_mode.hexa (2026-04-23).
  Fires every 30 min (1800s) — 3min avg continuous_scan 대비 idle 충분.
  Disabled by default. Install:
    cp config/launchd/com.hexa-lang.sleep_mode.plist \
       ~/Library/LaunchAgents/com.hexa-lang.sleep_mode.plist
    launchctl load ~/Library/LaunchAgents/com.hexa-lang.sleep_mode.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hexa-lang.sleep_mode</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ghost/core/hexa-lang/hexa</string>
        <string>/Users/ghost/core/hexa-lang/tool/sleep_mode.hexa</string>
        <string>--commit</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ghost/core/hexa-lang</string>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>StandardOutPath</key>
    <string>/tmp/hexa_sleep_mode.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/hexa_sleep_mode.err</string>

    <key>Disabled</key>
    <false/>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
