# config/launchd/com.hexa-lang.self_continuous_scan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.hexa-lang.self_continuous_scan — Phase 5 self-meta continuous runner.
  Spec: docs/hexa_lang_meta_evolution_proposal_20260422.md §Phase 5.

  Runs hx_meta continuous-scan (Phase 1 + 2 + 3 cascade) every 12 hours.
  Install:
    cp config/launchd/com.hexa-lang.self_continuous_scan.plist \
       ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) \
       ~/Library/LaunchAgents/com.hexa-lang.self_continuous_scan.plist
  Uninstall:
    launchctl bootout gui/$(id -u)/com.hexa-lang.self_continuous_scan

  STATUS: stub — enable after Phase 1-6 tools ship with real logic.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hexa-lang.self_continuous_scan</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd "$HOME/core/hexa-lang" &amp;&amp; ./bin/hx_meta continuous-scan 2&gt;&amp;1 | tee -a state/hx_continuous_scan.log</string>
    </array>

    <key>StartInterval</key>
    <integer>43200</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/com.hexa-lang.self_continuous_scan.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/com.hexa-lang.self_continuous_scan.err</string>

    <key>Disabled</key>
    <true/>
</dict>
</plist>
