# launchd/com.hexa-lang.relock.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
        "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <!--
    com.hexa-lang.relock — OPTIONAL launchd agent for raw#1 P0d.

    The primary path is hx_unlock.hexa spawning a detached
    tool/relock_daemon.hexa child per unlock (session-local, no
    installation). This plist is provided for operators who prefer
    a system-registered watcher that polls every 60s and re-locks
    whenever an unlock has been outstanding for longer than its
    configured window.

    INSTALL (manual, opt-in):
      cp launchd/com.hexa-lang.relock.plist ~/Library/LaunchAgents/
      # edit WorkingDirectory below to match your hexa-lang checkout
      launchctl load ~/Library/LaunchAgents/com.hexa-lang.relock.plist

    UNINSTALL:
      launchctl unload ~/Library/LaunchAgents/com.hexa-lang.relock.plist
      rm ~/Library/LaunchAgents/com.hexa-lang.relock.plist

    The spawned-per-unlock path does NOT need this file.
  -->
  <key>Label</key>
  <string>com.hexa-lang.relock</string>

  <key>ProgramArguments</key>
  <array>
    <string>./build/hexa_interp</string>
    <string>tool/relock_daemon.hexa</string>
    <string>--window</string>
    <string>300</string>
  </array>

  <!-- EDIT THIS to match your local checkout path -->
  <key>WorkingDirectory</key>
  <string>/Users/ghost/core/hexa-lang</string>

  <!-- Poll once per minute; each run is cheap (exits 0 if already locked). -->
  <key>StartInterval</key>
  <integer>60</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/hexa-lang.relock.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/hexa-lang.relock.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
