# launchd/dev.hexa-lang.foreign-ext-sweep.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.hexa-lang.foreign-ext-sweep</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/ghost/.hx/bin/hexa</string>
    <string>/Users/ghost/core/hexa-lang/tool/foreign_ext_sweep.hexa</string>
    <string>--sweep</string>
    <string>/Users/ghost/core</string>
  </array>
  <key>StartInterval</key>
  <integer>30</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/ghost/.hx/log/foreign-ext-sweep.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/ghost/.hx/log/foreign-ext-sweep.err</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>WS</key>
    <string>/Users/ghost/core</string>
    <key>HEXA_LANG</key>
    <string>/Users/ghost/core/hexa-lang</string>
    <key>HEXA_RESOLVER_NO_REROUTE</key>
    <string>1</string>
    <key>LC_ALL</key>
    <string>C</string>
  </dict>
</dict>
</plist>
